# supervisor_bloghook.ini
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root
autostart=true
autorestart=true

[program:ginweb]
command=/app/main
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/ginweb_stdout.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=3
redirect_stderr=true